# debian/usr/lib/systemd/system/monad-cruft.timer
[Timer]
OnCalendar=hourly
Unit=monad-cruft.service

[Install]
WantedBy=multi-user.target
